# vps-agent/systemd/mt5-agent.service
[Unit]
Description=Extrack MT5 VPS Agent
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/mt5-agent
EnvironmentFile=/home/ubuntu/mt5-agent/.env
Environment=AGENT_SERVICE_MODE=1
Environment=RUN_ONCE=0
ExecStart=/home/ubuntu/mt5-agent/agent.sh
Restart=always
RestartSec=10
KillMode=process
StandardOutput=append:/home/ubuntu/mt5-agent/logs/systemd.out.log
StandardError=append:/home/ubuntu/mt5-agent/logs/systemd.err.log

[Install]
WantedBy=multi-user.target
